# scripts/dist/services/photoprism/run
#!/command/execlineb -P
with-contenv
/scripts/cmd.sh /opt/photoprism/bin/photoprism start
